# dot_config/systemd/user/hyprsunset-check.service
[Unit]
Description=Check and apply nightlight state

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprsunset-apply
