# dot_config/systemd/user/hyprsunset-check.timer
[Unit]
Description=Periodic nightlight check

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
