# dot_config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate wallpaper randomly

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-random
